# supervisord.conf
[supervisord]
nodaemon=true


[program:kibana]
command = /opt/kibana/bin/kibana
redirect_stderr=true
