# supervisord.conf
[supervisord]
nodaemon=true


[program:kibana]
command = /opt/kibana-$KIBANA_VERSION/bin/kibana
redirect_stderr=true
stdout_logfile=/var/log/supervisor/kibana.log
